# watchlist-monitor.service
[Unit]
Description=Longbridge Watchlist Real-Time Monitor
After=network.target

[Service]
# User running the service (change if needed)
User=appuser
WorkingDirectory=/opt/longport_stock_monitor
# Use the python executable from the virtual environment
ExecStart=/opt/longport_stock_monitor/venv/bin/python3 -m src.monitor.watchlist_monitor
Restart=always
RestartSec=5
# Output logs to file
StandardOutput=append:/opt/longport_stock_monitor/logs/monitor.log
StandardError=append:/opt/longport_stock_monitor/logs/error.log
# Environment variables if needed
# Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
